# scripts/alc-backup.cron.example
# Runs the ALC Client Inquiry backup every day at 1:30 AM.
# Edit the project path if your VPS stores the app somewhere else.
#
# Install with:
#   crontab -e
#
# Then paste the non-comment line below:

30 1 * * * cd /home/agusanlending/ALC-client-inquiry-system && mkdir -p backups && APP_DIR=/home/agusanlending/ALC-client-inquiry-system RETENTION_DAYS=14 bash scripts/backup.sh >> backups/backup.log 2>&1
